# tools/linux/agent/systemd/appsandbox-clipboard.service
[Unit]
Description=App Sandbox clipboard daemon (XCB ↔ vsock :5/:6)
Documentation=https://github.com/jamesstringer90/appsandbox
After=graphical-session.target
PartOf=graphical-session.target

[Service]
# User-level unit. Inherits DISPLAY / WAYLAND_DISPLAY / XAUTHORITY from
# the session environment via systemd's import-environment hook on
# graphical-session.target, so libxcb finds XWayland without us setting
# them explicitly here.
Type=simple
ExecStart=/usr/local/bin/appsandbox-clipboard
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
