# tools/linux/agent/systemd/appsandbox-input.service
[Unit]
Description=App Sandbox input daemon (uinput injection ← vsock :3)
Documentation=https://github.com/jamesstringer90/appsandbox
After=systemd-udevd.service appsandbox-agent.service

[Service]
Type=simple
ExecStart=/usr/local/bin/appsandbox-input
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
# /dev/uinput requires CAP_SYS_ADMIN by default; root is simplest.
User=root

[Install]
WantedBy=multi-user.target
