# tools/linux/agent/systemd/appsandbox-firstboot.service
[Unit]
Description=App Sandbox first-boot setup (desktop install, kernel modules)
Documentation=https://github.com/jamesstringer90/appsandbox
# Wait for the agent to be up (it applies the host-assigned static IP for
# NAT-mode VMs via netplan). Then wait for network-online so apt has
# internet before we try to download ubuntu-desktop-minimal.
Requires=appsandbox-agent.service
After=appsandbox-agent.service
Wants=network-online.target
After=network-online.target
# Idempotent guard: cloud-init touches this marker once at install time;
# ExecStartPost removes it. Subsequent boots skip the service entirely.
ConditionPathExists=/etc/appsandbox-firstboot.marker

[Service]
Type=oneshot
ExecStart=/opt/appsandbox/setup.sh
ExecStartPost=/bin/rm -f /etc/appsandbox-firstboot.marker
# apt install of ubuntu-desktop-minimal can take 10-15 min on a slow link.
TimeoutStartSec=45min
StandardOutput=journal
StandardError=journal
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
