# tools/linux/agent/systemd/appsandbox-audio.service
[Unit]
Description=App Sandbox audio daemon (snd-aloop → vsock :4)
Documentation=https://github.com/jamesstringer90/appsandbox
# snd-aloop must be loaded before we open hw:Loopback,1,0.
After=systemd-modules-load.service appsandbox-agent.service
Requires=systemd-modules-load.service

[Service]
Type=simple
ExecStart=/usr/local/bin/appsandbox-audio
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
# libasound opens /dev/snd/* which is in the `audio` group; running as
# root is simplest and matches the other in-VM daemons.
User=root

[Install]
WantedBy=multi-user.target
