# tools/linux/agent/systemd/appsandbox-agent.service
[Unit]
Description=App Sandbox guest agent
Documentation=https://github.com/jamesstringer90/appsandbox
# Agent talks to the host over vsock (hv_sock), not IP. Depending on
# network-online.target would gate agent startup behind cloud-init's
# slow init stage and any cidata network ordering — for no benefit,
# since hv_sock is up as soon as the vmbus driver is loaded.
After=local-fs.target
Wants=local-fs.target

[Service]
Type=simple
ExecStart=/usr/local/bin/appsandbox-agent
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
# Need CAP_NET_ADMIN to run `netplan apply` for set_ip, and we shell out
# to systemctl poweroff / reboot. Run as root — keeps things simple and
# matches the Windows agent which runs as LocalSystem.
User=root

[Install]
WantedBy=multi-user.target
