# tools/linux/agent/systemd/appsandbox-display.service
[Unit]
Description=App Sandbox display daemon (KMS framebuffer capture → vsock :2)
Documentation=https://github.com/jamesstringer90/appsandbox
After=appsandbox-agent.service
# No After=graphical.target — that introduces an ordering cycle with
# multi-user.target on Desktop installs (multi-user→ppd→graphical→us).
# The daemon polls /dev/dri/cardN itself and retries until the active
# plane fb_id is non-zero, so it tolerates starting before GDM.

[Service]
Type=simple
ExecStart=/usr/local/bin/appsandbox-display
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
# Needs CAP_SYS_ADMIN for drmModeGetFB2 (kernel hardening since 4.15).
# Running as root keeps that simple and also gives us /dev/dri/cardN access
# without juggling the render group.
User=root

[Install]
WantedBy=multi-user.target
