# tools/linux/asb_drm/systemd-asb-evict-simpledrm.service
[Unit]
Description=Unbind simpledrm so asb_drm is the only display GNOME sees
Documentation=https://github.com/jamesstringer90/appsandbox
DefaultDependencies=no
After=systemd-modules-load.service
Before=display-manager.service graphical.target
ConditionPathExists=/sys/bus/platform/drivers/simple-framebuffer

[Service]
# simpledrm is built into the Ubuntu kernel (CONFIG_DRM_SIMPLEDRM=y) so a
# modprobe blacklist is a no-op. The aperture_remove kernel API didn't
# evict it either on this kernel because simpledrm's binding doesn't go
# through the aperture registry the way real-GPU drivers' do.
#
# What does work: unbind the simple-framebuffer driver from its platform
# device. After this, /dev/dri/card0 disappears and GNOME stops seeing a
# second monitor next to asb_drm.
Type=oneshot
ExecStart=/bin/sh -c 'for d in /sys/bus/platform/drivers/simple-framebuffer/simple-framebuffer.*; do [ -e "$d" ] && basename "$d" > /sys/bus/platform/drivers/simple-framebuffer/unbind; done; true'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
